# deploy/wfilemanager-heartbeat.timer
[Unit]
Description=Send the wFileManager managed Supabase heartbeat twice daily

[Timer]
OnBootSec=2min
OnUnitActiveSec=12h
RandomizedDelaySec=15min
Persistent=true
Unit=wfilemanager-heartbeat.service

[Install]
WantedBy=timers.target
